# roles/deploy_opa_postgres/templates/container-postgresql.service.j2
# container-postgresql.service
# autogenerated by Podman 5.6.0
# Wed Apr 29 13:57:14 EDT 2026

[Unit]
Description=Podman container-postgresql.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Environment=POSTGRESQL_PASSWORD={{ pg_password }}
Restart=always
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--sdnotify=conmon \
	--replace \
	-d \
	--name postgresql \
	-p 5432:5432 \
	-e POSTGRESQL_PASSWORD \
	-e POSTGRESQL_USER=postgres \
	-e POSTGRESQL_DATABASE=compliance \
	-v compliance_pgdata:/var/lib/pgsql/data:Z registry.redhat.io/rhel9/postgresql-15
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
